# recipes-core/systemd/files/95-gpio-buttons.rules
SUBSYSTEM=="button",ENV{ACTION}=="released", RUN+="/bin/sh /lib/wifi/wps-hotplug.sh"
SUBSYSTEM=="button",ENV{ACTION}=="released", RUN+="/bin/sh /lib/wifi/wps-reset-hotplug.sh"
